# Python/dh-scheduler_0.conf
# MySQL Service

description     "dh-schedurel_0.py"
author          "Franco Parenti"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

# The default of 5 seconds is too low for mysql which needs to flush buffers
kill timeout 300

chdir /var/dhproc
exec python3 /var/dhproc/scheduler.py 0
respawn

